# init/vappio_ec2.conf
description     "EC2 startup"

start on (starting vappio_detect and vappio_ec2)

task

console output

exec run-parts /etc/vappio/amazonec2/
