# init/vappio_vmware.conf
description     "VMware startup"

start on (starting vappio_detect and vappio_vmware)

task

console output

exec run-parts /etc/vappio/vmware/
