# img-conf/etc/cron.d/ntpdate
# Need to run this every 5 minutes
*/5 * * * * /etc/init.d/ntpdate.sh
