# img-conf/etc/cron.d/ntpdate
*/5 * * * * root bash /etc/init.d/ntpdate.sh
