# img-conf/etc/update-motd.d/00-header
#!/bin/sh

uname -a
printf "%s\n" "CloVR VM $(lsb_release -s -d)"
