# img-conf/etc/update-motd.d/10-help-text
#!/bin/bash

if [ -d "/etc/vappio" ]
then
    cat /etc/vappio/issue
fi