# img-conf/etc/update-motd.d/10-help-text
#!/bin/bash

if [ -d "/etc/vappio" ]
then
    #Reuse pre-login text here
    cat /etc/issue
fi