# img-conf/etc/update-motd.d/50-ganglia-stats
#!/bin/bash


#which gstats > /dev/null
#if [ $? != 0 ]
#then
#	gstats | grep -v gexec
#fi
#grep top info here: load, cpu, mem, 

